# comparison/ansible/roles/mellanox_nic/templates/rc.local.j2
#!/bin/sh -e
echo 2 > /sys/class/net/{{ mellanox_pf }}/device/sriov_numvfs
