# comparison/ansible/roles/multus_cni/files/kubelet-watcher.path
[Path]
PathModified=/etc/cni/net.d

[Install]
WantedBy=multi-user.target
