# comparison/ansible/roles/multus_cni/files/rc.local
#! /bin/bash

/bin/bash /opt/cnf_testbed/vf_setup.sh

exit 0
