# comparison/ansible/roles/vfs_vswitch/templates/rc.local.j2
#!/bin/sh -e
echo 8 > /sys/class/net/{{ interface }}/device/sriov_numvfs
ip link set {{ interface }} vf 0 trust on
ip link set {{ interface }} vf 1 trust on
